# dev_scripts/ansible/roles/common/templates/pia.service
[Unit]
Description=pia-tunnel autostart
After=network.target auditd.service sshd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pia-boot.sh

[Install]
WantedBy=multi-user.target